# supervisor.conf
[supervisord]
nodaemon=true

[program:gearmand]
command=/usr/sbin/gearmand

[program:gearman-ui]
command=/usr/bin/php -S 0.0.0.0:4731 -t /gearman-ui/web /gearman-ui/web/index.php
